# crontab
*/15  * * * * /opt/idle-sleep/ecosleep.sh auto
