# profile.d/ecosleep.sh
#!/bin/bash

alias ecosleep="/bin/sudo /opt/idle-sleep/ecosleep.sh"
